# bin/greent-graphql.service
[Unit]
Description=The nginx HTTP and reverse proxy server
After=network.target remote-fs.target nss-lookup.target

[Service]
User=evryscope
ExecStart=/home/evryscope/dev/greent/bin/systemd-start
ExecStop=/home/evryscope/dev/greent/bin/systemd-stop

[Install]
WantedBy=multi-user.target
